# configs/etc/profile.d/ls.sh
alias ls='LC_ALL=C ls -1Ahv --color --show-control-chars --quoting-style=shell --group-directories-first'

